# systemd/time-logger.service
[Unit]
Description=TIme Logger Service

[Service]
ExecStart=/home/ubuntu/linux-lab-notes/systemd/time-logger.sh
